# scripts/crontab
*/5 * * * * /usr/local/bin/python ../etl.py >> /var/log/cron.log 2>&1